# systemd/personal-xp.service
[Unit]
Description=personal-xp FastAPI app
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/personal-xp
EnvironmentFile=/root/.wallet/personal-xp.env
ExecStart=/opt/personal-xp/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8004 --workers 2 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
